# openclaw-config/launchagents/ai.openclaw.auto-doc.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>ai.openclaw.auto-doc</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/alfredpifi/clawd/workspace-roberto/scripts/auto-doc.sh</string>
	</array>
	<key>StartInterval</key>
	<integer>900</integer>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/tmp/auto-doc.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/auto-doc.log</string>
</dict>
</plist>
